# cnc-dashboard.service
[Unit]
Description=CNC Dashboard Startup Service
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=cv
WorkingDirectory=/home/cv/cnc-dashboard
ExecStart=/home/cv/cnc-dashboard/startup.sh
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/cnc-dashboard.log
StandardError=append:/var/log/cnc-dashboard.error.log

[Install]
WantedBy=multi-user.target 